# ts_entities/rootfs/etc/services.d/ts_entities/run
#!/usr/bin/with-contenv sh
echo "Starting TS Entities add-on..."
echo "Node: $(node --version)"
exec node /app/packages/addon/dist/index.js
